# cron/cron.daily/backup
#!/bin/sh
gsutil rsync -d -r /var/backups/borg $GCS_BACKUP_URI
